# LaunchAgents/com.verknowsys.SyncEmAll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>Launch SyncEmAll</string>
	<key>Program</key>
	<string>/Users/Shared/Small/bin/small_wrapper</string>
	<key>EnvironmentVariables</key>
	  <dict>
	    <key>PATH</key>
	    <string>/usr/bin:/bin</string>
	    <key>HOME</key>
	    <string>/Users/dmilith</string>
	  </dict>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>

<!-- @@@@LingonWhatStart:/Software/Zsh/exports/zsh -c /usr/local/bin/small@@@@LingonWhatEnd -->
